# src/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
loglevel=info

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; Icecast first — Liquidsoap needs it up before it can connect
[program:icecast]
command=icecast2 -c /etc/radio/icecast.xml
autostart=true
autorestart=true
priority=10
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Python service second — Liquidsoap pushes track-change events to it
[program:radio-service]
command=uvicorn app.main:app --host 127.0.0.1 --port 8081
directory=/app
autostart=true
autorestart=true
priority=20
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Liquidsoap last — will retry if Icecast isn't ready yet
[program:liquidsoap]
command=liquidsoap /etc/radio/radio.liq
autostart=true
autorestart=true
priority=30
startretries=10
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
